# table-manager.service
# Table Manager systemd service
#
# Установка:
#   1. Скопировать: sudo cp table-manager.service /etc/systemd/system/
#   2. Заменить USER_PATH на домашнюю директорию пользователя
#   3. Заменить PROJECT_PATH на путь к проекту
#   4. Перезагрузить: sudo systemctl daemon-reload
#   5. Включить: sudo systemctl enable --now table-manager
#
# Проверка:
#   sudo systemctl status table-manager
#   journalctl -u table-manager -f

[Unit]
Description=Table Manager Next.js Application
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=qwen
Group=qwen
WorkingDirectory=/home/qwen/table-manager
EnvironmentFile=/home/qwen/table-manager/.env
Environment=NODE_ENV=production
Environment=PATH=/home/qwen/.nvm/versions/node/v20.20.2/bin:/usr/local/bin:/usr/bin:/bin

ExecStart=/home/qwen/.nvm/versions/node/v20.20.2/bin/npm start
# Next.js не поддерживает graceful reload — использовать systemctl restart

Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=table-manager

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# Resource limits (Node.js нужен много file descriptors)
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
